# deployments/omatab.service
[Unit]
Description="Omatab service"

[Service]
Type=simple
ExecStart=/home/omatab/start.sh
#ExecStart=node -v
Restart=always
User=omatab

[Install]
WantedBy=multi-user.target
